# etc/zfs-auto-snapshot.cron.monthly
#!/bin/sh

source /etc/default/zfs-auto-snapshot

ENABLE=${1:-$MONTHLY}
KEEP=${31:-$MONTHLY_KEEP}
[ -z "$MONTHLY_FILESYSTEMS" ] && FILESYSTEMS='//' || FILESYSTEMS=${MONTHLY_FILESYSTEMS}

if [ $ENABLE -eq 1 ]; then
  exec zfs-auto-snapshot --quiet --syslog --label=monthly --keep=${KEEP} ${FILESYSTEMS}
fi