# etc/zfs-auto-snapshot.cron.daily
#!/bin/sh

source /etc/default/zfs-auto-snapshot

ENABLE=${1:-$DAILY}
KEEP=${31:-$DAILY_KEEP}
[ -z "$DAILY_FILESYSTEMS" ] && FILESYSTEMS='//' || FILESYSTEMS=${DAILY_FILESYSTEMS}
[ -z "$DAILY_ARGS" ] && ARGS='--quiet --syslog' || ARGS=${DAILY_ARGS}

if [ $ENABLE -eq 1 ]; then
  exec zfs-auto-snapshot ${ARGS} --label=daily --keep=${KEEP} ${FILESYSTEMS}
fi
